# docker/supervisord/rr-server.ini
[program:rr-server]
process_name=%(program_name)s_%(process_num)02d
command=rr -c /app/.rr.yaml -w /app --dotenv /app/.env serve
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=1
startsecs=0
stopwaitsecs=3600
stdout_logfile=/var/log/supervisor/roadrunner_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/roadrunner_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
